# systemd-services/whale-tracker-commands.service
[Unit]
Description=Hyperliquid Whale Tracker - Telegram Command Handler
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/whale-tracker
Environment=PATH=/usr/bin:/usr/local/bin
ExecStart=/usr/bin/python3 telegram_command_handler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 